# cronie/config.crontab
*/5 * * * * /usr/local/bin/mailsync

# vim:ft=crontab
