# packages/jack/jack.service
[Unit]
Description=JACK service

[Service]
ExecStart=bash start.sh
ExecStartPost=bash connect.sh
ExecStopPost=bash stop.sh
WorkingDirectory=%h/.dotfiles/packages/jack
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
