# packages/compton/compton.service
[Unit]
Description=Compton service

[Service]
ExecStart=/usr/bin/compton --config config
WorkingDirectory=%h/.dotfiles/packages/compton
Environment=DISPLAY=:0
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
