# audio/jack/jack.service
[Unit]
Description=JACK service

[Service]
ExecStart=/bin/bash start.sh
ExecStartPost=/bin/bash connect.sh
ExecStopPost=/bin/bash stop.sh
WorkingDirectory=%h/.dotfiles/audio/jack

[Install]
WantedBy=multi-user.target
