# packages/mail-tunnel/mail-tunnel.service
[Unit]
Description=Mail tunnel service

[Service]
ExecStart=/home/soywod/.dotfiles/scripts/mail-tunnel.sh

[Install]
WantedBy=multi-user.target
