# display/compton/compton.service
[Unit]
Description=Compton service

[Service]
ExecStart=/usr/bin/compton --config config
WorkingDirectory=%h/.dotfiles/display/compton
Environment=DISPLAY=:0
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
