# packages/dunst/dunst.service
[Unit]
Description=Dunst service

[Service]
ExecStart=/usr/bin/dunst
WorkingDirectory=%h/.dotfiles/packages/dunst
Environment=DISPLAY=:0
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
